# Kindle/etc/upstart/startup.conf
description "Startup Tasks"

start on framework_ready

script
  # Start weather updater
  /mnt/us/local/bin/weather &
end script
